# genesisrpc.service
[Unit]
Description=genesisL1 node
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/root/go/bin/evmosd start --log_level info
Restart=on-failure
RestartSec=60 
LimitNOFILE=50000 

[Install]
WantedBy=multi-user.target
